# packages/sbn/sbn-77.tar.gz/sbn-77/files/sbn.service
[Unit]
Description=Skull, Bones and Number (OTP-CR-117/19)
Requires=network.target
After=network.target

[Service]
Type=simple
User=bart
Group=bart
WorkingDirectory=/home/bart/.sbn
ExecStart=/home/bart/.local/pipx/venvs/sbn/bin/sbnd
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
